# net-dns/dnsexit-ipUpdate/files/dnsexit-ipUpdate.service
[Unit]
Description=dnsexit Dynamic DNS Update Client
After=network.target

[Service]
Type=forking
ExecStart=EROOTusr/sbin/dnsexit-ipUpdate

[Install]
WantedBy=multi-user.target
 
